# meta-common/recipes-common/test-app/files/test-app.service
[Unit]
Description=Test App/Service

[Service]
ExecStart=/usr/bin/test-app
Restart=on-failure

[Install]
WantedBy=multi-user.target
